# files/ycsb-install.service
[Unit]
Description=Create the root user

[Service]
Type=simple
User=root

ExecStart=/usr/bin/ycsb-install

[Install]
WantedBy=multi-user.target
